# supervisor.d/dev2_outloud.ini
[program:uwsgi_dev2_outloud]
command=/usr/local/bin/uwsgi --ini "/usr/local/etc/uwsgi/dev2_outloud.ini" --uid outloud
stdout_logfile=/var/log/uwsgi/dev2_outloud_supervisor.log
stderr_logfile=/var/log/uwsgi/dev2_outloud_err.log
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs = 60
stopsignal=INT
